# systemd/goals-webapp.service
[Unit]
Description=Goals WebApp API (FastAPI)
After=network.target

[Service]
Type=simple
# Замените на своего пользователя и реальные пути (НЕ оставляйте /path/to/...)
User=root
WorkingDirectory=/root/shaolen
Environment=PATH=/usr/bin:/usr/local/bin
EnvironmentFile=/root/shaolen/.env
# Если используете venv: /root/shaolen/venv/bin/python3
ExecStart=/root/shaolen/venv/bin/python3 /root/shaolen/webapp_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
