# systemd/goals-bot.service
[Unit]
Description=Goals Telegram Bot (Shaolen)
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/telegram_goals_bot
Environment=PATH=/usr/bin:/usr/local/bin
EnvironmentFile=/path/to/telegram_goals_bot/.env
ExecStart=/usr/bin/python3 /path/to/telegram_goals_bot/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
